# config/supervisor-consul-template.conf
[program:consul-template]
command=/usr/local/bin/consul-template -config=/etc/consul-template/config.hcl
autorestart=true
priority=1
