# config/supervisor-rsyslog.conf
[program:rsyslog]
command=rsyslogd -n -c3
autorestart=true
priority=2
