# config/supervisor-postfix.conf
[program:postfix]
command=/usr/lib/postfix/master -d
autorestart=true
priority=3
